# files/profile.d/set-hadoop-prefix.sh
#!/bin/bash

export HADOOP_PREFIX=/usr/local/hadoop
export HADOOP_CONF_DIR=/etc/hadoop
export HADOOP_LOG_DIR=/var/local/hadoop/logs

PATH=${PATH}:/usr/local/hadoop/bin:/usr/local/hadoop/sbin
