# files/profile.d/set-spark-prefix.sh
#!/bin/bash

export SPARK_HOME=/usr/local/spark

PATH=${PATH}:/usr/local/spark/bin
